# config/udp-custom.service
[Unit]
Description=UDP Custom Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/udp
ExecStart=/usr/local/bin/udp-custom server --config /etc/udp/config.json
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target